# otbr/otbr-agent.initd
#!/sbin/openrc-run

name="otbr-agent"
description="OpenThread Border Router Agent"

command="/usr/sbin/otbr-agent"
command_args="$OTBR_AGENT_OPTS"
command_user="${OTBR_USER:-root}"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
	need net avahi-daemon
	after firewall
}

_detect_infra_if() {
	# Auto-detect backbone interface: prefer eth*, then wlan*, then any other
	local iface
	for iface in /sys/class/net/eth* /sys/class/net/en* /sys/class/net/wlan* /sys/class/net/*; do
		[ -d "$iface" ] || continue
		iface=$(basename "$iface")
		# Skip loopback, Thread interface, and virtual interfaces
		case "$iface" in
			lo|wpan*|docker*|veth*|br-*|tun*) continue ;;
		esac
		# Check if interface is UP and has an IPv4 address
		if ip link show "$iface" 2>/dev/null | grep -q "state UP" && \
		   ip -4 addr show "$iface" 2>/dev/null | grep -q "inet "; then
			echo "$iface"
			return 0
		fi
	done
	return 1
}

start_pre() {
	# Auto-detect backbone interface if set to "auto"
	if [ "$OTBR_INFRA_IF" = "auto" ]; then
		OTBR_INFRA_IF=$(_detect_infra_if) || {
			eerror "Could not auto-detect backbone interface"
			eerror "Set OTBR_INFRA_IF in /etc/conf.d/otbr-agent"
			return 1
		}
		einfo "Auto-detected backbone interface: $OTBR_INFRA_IF"
		# Rebuild command args with detected interface
		command_args="-I $OTBR_THREAD_IF -B $OTBR_INFRA_IF --rest-listen-address $OTBR_REST_ADDR --rest-listen-port $OTBR_REST_PORT spinel+hdlc+uart://$OTBR_RADIO_DEVICE"
	fi

	# Load TUN module if not already loaded
	if [ ! -e /dev/net/tun ]; then
		modprobe tun || {
			eerror "Failed to load TUN module"
			return 1
		}
	fi

	# Load modules for firewall/NAT64 (optional, non-fatal)
	modprobe ip6table_filter 2>/dev/null || true
	modprobe ip6_tables 2>/dev/null || true

	# Enable IP forwarding for border routing
	if [ -f /etc/sysctl.d/otbr.conf ]; then
		sysctl -p /etc/sysctl.d/otbr.conf >/dev/null 2>&1
	else
		# Fallback if sysctl file missing
		sysctl -w net.ipv6.conf.all.forwarding=1 >/dev/null 2>&1
		sysctl -w net.ipv4.ip_forward=1 >/dev/null 2>&1
	fi

	# Ensure radio device exists
	if [ -n "$OTBR_RADIO_DEVICE" ] && [ ! -c "$OTBR_RADIO_DEVICE" ]; then
		eerror "Radio device $OTBR_RADIO_DEVICE not found"
		return 1
	fi

	checkpath --directory --owner "$command_user" /var/lib/otbr
}
